# unix/profile.d/01_ulimit
#!/bin/sh

# Core size
ulimit -c 500000

