# unix/profile.d/02_prompt
#!/bin/sh

export PS1='${USER}@$(hostname):${PWD/$HOME/\~} $ '

