# unix/profile.d/11_inferno
#!/bin/sh

# inferno wm args
export EMU="-r$HOME/src/inferno -c0 -g1024x768"

