# unix/profile.d/12_luarocks
#!/bin/sh

export PATH=$PATH:$HOME/.luarocks/bin
if which luarocks >/dev/null ; then
	eval `luarocks path`
fi

