# unix/profile.d/21_htop
#!/bin/sh

export HTOPRC=$HOME/.htoprc

