# unix/profile.d/22_kubectl
#!/bin/sh

if command -v kubectl>/dev/null && test -n "${BASH}"; then
	source <(kubectl completion bash)
fi

