# unix/profile.d/30_noa11y
#!/bin/sh

export NO_AT_BRIDGE=1
export GTK_A11Y=none

